# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:lighttpd]
command=/bin/bash -c "exec env LANG=C /usr/bin/mrtg /mrtgcfg/mrtg.cfg restart"